# raspirecord.service
[Unit]
Description=RaspiRecordSync Auto Recording Service
After=network.target

[Service]
Type=simple
User=tspol
WorkingDirectory=/home/tspol/RaspiRecordSync
ExecStart=/bin/bash /home/tspol/RaspiRecordSync/start_realtime_system.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONPATH=/home/tspol/RaspiRecordSync

[Install]
WantedBy=multi-user.target 